# playground/day14-server-health-check-script/systemd/health-moniter.timer
[Unit]
Description=Run health monitor every hour

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
